# build/debian/etc/cron.d/groupoffice-cron
* * * * * www-data php /usr/share/groupoffice/cron.php /etc/groupoffice/config.php

