# debian-groupoffice-com/etc/cron.d/groupoffice-cron
* * * * * www-data php /usr/share/groupoffice/groupofficecli.php -r=core/cron/run -c=/etc/groupoffice/config.php -q > /dev/null
