# debian-groupoffice-com/etc/cron.d/groupoffice-log
0 0 * * * root /usr/bin/php /usr/share/groupoffice/groupofficecli.php -r=log/log/rotate -c=/etc/groupoffice/config.php > /dev/null
