# debian-groupoffice-mailserver/etc/cron.d/groupoffice-mailserver
#Updates usage information in the servermanager.
1 2 * * * root /usr/share/groupoffice/groupofficecli.php -r=postfixadmin/mailbox/cacheUsage -c=/etc/groupoffice/config.php -q > /dev/null
