# build/debian/debian/postinst
#!/bin/sh
# source debconf stuff
. /usr/share/debconf/confmodule
dbc_generate_include_owner="root:www-data"
dbc_generate_include_perms="0640"
dbc_generate_include=template:/etc/groupoffice/config.php
dbc_generate_include_args="-U -o template_infile='/usr/share/groupoffice/config.php.tpl'"
       
. /usr/share/dbconfig-common/dpkg/postinst
dbc_go groupoffice "$@"

#This approach caused an error when reloading while apache was updated too. It was not running.
#a2enconf groupoffice
#invoke-rc.d apache2 reload

if [ -f "/usr/share/apache2/apache2-maintscript-helper" ]; then
  . /usr/share/apache2/apache2-maintscript-helper
  apache2_invoke enconf groupoffice
fi

chown www-data:www-data /var/lib/groupoffice

if [ -z "$2" ]; then
  echo "Please finish the install in the browser";
else
  /usr/share/groupoffice/cli.php core/System/upgrade
fi

